# root/etc/s6-overlay/s6-rc.d/svc-piper/run
#!/command/with-contenv bash
# shellcheck shell=bash

NVIDIA_PACKAGE_DIR="/lsiopy/lib/python3.12/site-packages/nvidia"

for dir in "${NVIDIA_PACKAGE_DIR}"/*; do
    if [ -d "$dir/lib" ]; then
        export LD_LIBRARY_PATH="$dir/lib:$LD_LIBRARY_PATH"
    fi
done

unset UPDATE_MODELS

if [[ -z ${LOCAL_ONLY} ]]; then
    UPDATE_MODELS=true
fi


LD_LIBRARY_PATH="/lsiopy/lib/python3.12/site-packages/onnxruntime/capi:$LD_LIBRARY_PATH"

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 10200" \
        s6-setuidgid abc python3 -m wyoming_piper \
        --piper '/lsiopy/bin/piper' \
        --uri 'tcp://0.0.0.0:10200' \
        --length-scale "${PIPER_LENGTH:-1.0}" \
        --noise-scale "${PIPER_NOISE:-0.667}" \
        --noise-w "${PIPER_NOISEW:-0.333}" \
        --speaker "${PIPER_SPEAKER:-0}" \
        --voice "${PIPER_VOICE}" \
        --max-piper-procs "${PIPER_PROCS:-1}" \
        --data-dir /config \
        --download-dir /config \
        --cuda \
        ${UPDATE_MODELS:+--update-voices} \
        ${STREAMING:+--streaming}

